# etc/systemd/system/ollama.service.d/aorus-egpu.conf
# Drop-in for ollama.service on AORUS RTX 5090 eGPU systems.
#
# The eGPU goes through a manual bind sequence at boot
# (aorus-egpu-compute-load-nvidia.service); ollama depends on it
# completing so the GPU is bound before CUDA discovery runs.
# Persistenced is a dependency for warmup-latency reasons (see
# docs/architecture.md Problem 2 — RECLASSIFIED 2026-05-08; saves
# ~1.3s GSP-boot tax per first-open after LAST-CLOSE).
#
# UPDATED 2026-05-08: removed dependency on
# aorus-egpu-uvm-keepalive.service. UVM-side close-path bug class
# was empirically demonstrated benign on the current driver stack
# (Patch 0030 + n=6 UVM probes; H22 ledger entry); uvm-keepalive
# retired the same day. With the dependency present here,
# `systemctl disable aorus-egpu-uvm-keepalive` did NOT actually
# prevent the service running — Requires= here was pulling it in
# whenever ollama started. Removing both After= and Requires= so
# the retirement is real.
#
# Requires= (not just After=) means ollama will not start if any
# remaining dependency failed — better to fail loudly than to run
# against a non-bound device.
#
# ConditionPathExists guards against running on a host where the
# GPU was not enumerated (e.g., Thunderbolt link down at boot).

[Unit]
After=aorus-egpu-compute-load-nvidia.service
After=nvidia-persistenced.service
Requires=aorus-egpu-compute-load-nvidia.service
Requires=nvidia-persistenced.service
ConditionPathExists=/sys/bus/pci/devices/0000:04:00.0
