# etc/systemd/system/nvidia-persistenced.service.d/aorus-egpu.conf
# Drop-in for nvidia-persistenced on AORUS RTX 5090 eGPU systems.
#
# Persistenced is load-bearing on this stack: it holds /dev/nvidia0 and
# /dev/nvidiactl open so that nvidia-smi (and any other NVML caller) is
# never the "first open after last close" of the device. Without this,
# the second nvidia-smi invocation in any module-load session hard-freezes
# the host (close-path teardown wedges the next open on this hardware).
#
# Ordering: must start AFTER aorus-egpu-compute-load-nvidia.service has
# bound the GPU to the nvidia driver. Requires= ensures persistenced will
# fail to start (rather than running against a non-bound device) if the
# bind step failed.
#
# Restart=no: do NOT auto-restart on failure. If persistenced dies while
# nvidia is loaded, restarting it would close+reopen device files and
# likely freeze the host. A dead daemon makes the system freeze-prone but
# does not actively make things worse; reboot to recover.

[Unit]
After=aorus-egpu-compute-load-nvidia.service
Requires=aorus-egpu-compute-load-nvidia.service
ConditionPathExists=/sys/bus/pci/devices/0000:04:00.0

[Service]
Restart=no
