# etc/udev/rules.d/79-aorus-egpu-no-autoload.rules
# Stop generic udev PCI modalias autoload from binding the eGPU before
# nvidia-persistenced is ready. NVIDIA is loaded explicitly through
# aorus-5090-compute-load-nvidia.service.
#
# Why both `driver_override` and clearing MODALIAS:
#   - driver_override blocks drivers_autoprobe from binding any registered driver
#     other than the named one (which won't ever exist).
#   - Clearing ENV{MODALIAS} prevents systemd-udevd's modalias matcher from
#     loading the nvidia kernel module by autoload at all.
#
# Without both, the freeze-prone "first NVIDIA probe before persistenced is up"
# can occur before our service has bound the GPU.

ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", ATTR{device}=="0x2b85", ATTR{driver_override}="aorus_5090_manual", ENV{MODALIAS}=""
ACTION=="add|bind", SUBSYSTEM=="pci", ATTR{vendor}=="0x10de", ATTR{device}=="0x22e8", ATTR{driver_override}="aorus_5090_disabled", ENV{MODALIAS}="", RUN+="/usr/local/sbin/aorus-5090-disable-audio %k"
