# etc/systemd/system/aorus-egpu-wpr2-recovery.service
[Unit]
Description=AORUS RTX 5090 eGPU WPR2-stuck recovery (Lever R Tier 1 v3)
Documentation=file:///root/aorus-5090-egpu/docs/lever-R-design.md
Documentation=file:///root/aorus-5090-egpu/docs/reliability-hypothesis-ledger.md

# Order: run AFTER the loader and its dependents have had a chance to settle.
# Use After= without Requires= because our PURPOSE is to recover from a
# failure of those services. If they succeeded, this helper is a no-op.
After=aorus-egpu-compute-load-nvidia.service
After=nvidia-persistenced.service
After=aorus-egpu-uvm-keepalive.service

# Wants= (not Requires=) so we run even if a prior service in the chain
# failed. The helper itself checks current state and decides whether to
# act.
Wants=aorus-egpu-compute-load-nvidia.service

# Skip cleanly if the eGPU isn't connected.
ConditionPathExists=/sys/bus/pci/devices/0000:04:00.0

[Service]
Type=oneshot
RemainAfterExit=yes

# The helper exits 0 on no-op-needed or successful recovery, exit
# codes 1-3 for various failure modes. The unit is "successful" only
# when exit=0; non-zero exits show as failed in `systemctl status`.
ExecStart=/usr/local/sbin/aorus-egpu-wpr2-recovery

StandardOutput=journal
StandardError=journal

# Don't restart on failure — recovery is one-shot. If it failed,
# manual intervention is required.
Restart=no

[Install]
WantedBy=multi-user.target
