# etc/systemd/system/aorus-egpu-uvm-keepalive.service
[Unit]
Description=AORUS RTX 5090 eGPU UVM keep-alive (close-path bug mitigation)
Documentation=file:///root/aorus-5090-gpu/docs/architecture.md
Documentation=file:///root/ollama/docs/freeze-2026-05-02-1032.md

# Order: only after the loader has bound nvidia AND pre-staged nvidia_uvm
# (which creates /dev/nvidia-uvm + /dev/nvidia-uvm-tools), and after
# persistenced is up. Combined with the WantedBy below, this places our open
# of UVM before any user-session CUDA service and after the platform's
# loader/persistenced chain. See /root/ollama/docs/freeze-2026-05-02-1032.md
# for the freeze evidence that motivated this unit.
After=aorus-5090-compute-load-nvidia.service
After=nvidia-persistenced.service
Requires=aorus-5090-compute-load-nvidia.service
Requires=nvidia-persistenced.service

# Skip cleanly if the loader didn't run / didn't pre-stage nvidia_uvm. With
# Requires= above this should be unreachable, but the explicit condition
# means the unit fails as 'condition' rather than 'failed' if invoked
# manually before the chain is up.
ConditionPathExists=/dev/nvidia-uvm
ConditionPathExists=/dev/nvidia-uvm-tools

[Service]
Type=exec
ExecStart=/usr/local/sbin/aorus-5090-uvm-keepalive

# Restart=no by design. If this holder dies, restart would mean close+reopen
# of the UVM fds, which is exactly the 1->0->1 transition we are preventing.
# A dead holder leaves the system in the same vulnerable state as no holder
# at all; recovery is reboot. Same reasoning as the persistenced drop-in.
Restart=no

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
