# etc/systemd/system/aorus-egpu-compute-load-nvidia.service
[Unit]
Description=Bind AORUS RTX 5090 eGPU to NVIDIA driver
Documentation=file:///root/aorus-5090-gpu/README.md

# Order: run after udev settles and Thunderbolt has authorized the eGPU,
# but before graphical.target so persistenced (which is After= this unit)
# is up before GDM starts probing devices.
After=systemd-udev-settle.service bolt.service
Wants=bolt.service systemd-udev-settle.service
Before=graphical.target

# Skip cleanly if the eGPU is not on PCI (disconnected, powered off, or
# Thunderbolt link not authorized in time). PCI path is stable for this
# topology - one TB controller, one device, address 0000:04:00.0.
ConditionPathExists=/sys/bus/pci/devices/0000:04:00.0

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/aorus-5090-compute-load-nvidia

[Install]
WantedBy=multi-user.target
